# systemd/todo-sync.service
[Unit]
Description=Sync todo list upon write

[Service]
type=oneshot
ExecStart=/home/arcwand/scripts/todo_sync.sh
