# systemd/todo-sync.path
[Path]
PathModified=/home/arcwand/Desktop/todo.txt

[Install]
WantedBy=default.target
